# crontab
*/5 * * * * /usr/bin/python3 /Users/chase/Code/live-bg/script.py
